# roles/websauna.site/templates/supervisor.conf
# Supervisor configuration to deploy funder on the production

[program:uwsgi]
{% if new_relic %}
command = {{ deploy_location }}/venv/bin/newrelic-admin run-program {{ uwsgi_config }}
{% else %}
command = {{deploy_location}}/venv/bin/uwsgi {{ uwsgi_config }}
{% endif %}
autostart=true
autorestart=true
directory={{ deploy_location }}
stdout_logfile={{ deploy_location }}/logs/uwsgi.log
redirect_stderr=true
stopsignal=QUIT
environment=NEW_RELIC_CONFIG_FILE={{deploy_location}}/{{ site_id }}-newrelic.ini 
user={{ wsgi_user }}
startsecs=15

{% if celery %}
[program:celery-beat]
command={{deploy_location}}/venv/bin/ws-celery {{ websauna_config_file }} -- beat --loglevel={{ celery_log_level }}
stderr_logfile={{ deploy_location }}/logs/celery-beat.stderr.log
stdout_logfile={{ deploy_location }}/logs/celery-beat.stdout.log
directory={{ deploy_location }}
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
user={{ wsgi_user }}

[program:celery-worker]
command={{deploy_location}}/venv/bin/ws-celery {{ websauna_config_file }} -- worker --loglevel={{ celery_log_level }}
stdout_logfile={{ deploy_location }}/logs/celery-worker.stdout.log
stderr_logfile={{ deploy_location }}/logs/celery-worker.stderr.log
directory={{ deploy_location }}
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
environment=C_FORCE_ROOT="true"
user={{ wsgi_user }}
{% endif %}

